# aerodrome-radar.service
[Unit]
Description=Aerodrome Radar Service
After=network.target

[Service]
Type=simple
User=yoshi
WorkingDirectory=/home/yoshi/backend
Environment=PYTHONIOENCODING=utf-8
Environment=LANG=en_US.UTF-8
ExecStart=/home/yoshi/backend/venv/bin/python /home/yoshi/backend/aerodrome_radar.py
Restart=always
RestartSec=10
StandardOutput=append:/home/yoshi/backend/radar.log
StandardError=append:/home/yoshi/backend/radar.log

[Install]
WantedBy=multi-user.target
